# systemd/site.service.example
# Example systemd service — generated by deploy-backend.sh
# Reference only. deploy-backend.sh installs the real one to /etc/systemd/system/
[Unit]
Description=Gunicorn for ${SITE_NAME}
After=network.target

[Service]
User=${GUNICORN_USER}
Group=${GUNICORN_GROUP}
WorkingDirectory=${WEBROOT}/${BACKEND_DIR}
EnvironmentFile=/etc/systemd/system/${SITE_NAME}.env
ExecStart=${PYTHON_BIN} -m gunicorn -k eventlet -w ${WORKERS} -b ${BACKEND_HOST}:${BACKEND_PORT} ${BACKEND_APP_SCRIPT}:app

[Install]
WantedBy=multi-user.target
